# apps/web/tripplanner/branches/dojo/service/tripplanner/run
#!/bin/sh
exec 2>&1
exec setuidgid bycycle /usr/local/bin/python2.5 -u ./serve.py

